# bin/drop-watcher.conf
[program:web_watcher]
command=python3 /home/shg/drop-watcher/agents/web_watcher.py
directory=/home/shg/drop-watcher
user=shg
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=DW_CONFIG_DIR="/etc/drop-watcher",DW_DATA_DIR="/var/lib/drop-watcher",DW_LOG_DIR="/var/log/drop-watcher",DW_ENV_FILE="/etc/drop-watcher/.env"
stdout_logfile=/var/log/drop-watcher/web_watcher.log
stderr_logfile=/var/log/drop-watcher/web_watcher_err.log

[program:watcher_signup]
command=/usr/local/bin/gunicorn -w 2 -b 127.0.0.1:5001 --graceful-timeout 5 --timeout 30 watcher_signup:app
directory=/home/shg/drop-watcher
user=shg
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10
environment=DW_CONFIG_DIR="/etc/drop-watcher",DW_DATA_DIR="/var/lib/drop-watcher",DW_LOG_DIR="/var/log/drop-watcher",DW_ENV_FILE="/etc/drop-watcher/.env"
stdout_logfile=/var/log/drop-watcher/watcher_signup.log
stderr_logfile=/var/log/drop-watcher/watcher_signup_err.log
